# etc/systemd/system/pressure-logger.service
[Unit]
Description=BMP180 logging service

[Service]
Type=simple
ExecStart=/usr/sbin/pressure-logger

[Install]
WantedBy=multi-user.target

